# www/zentaopms/bin/cron/sys.cron
# system cron.
#min   hour day month week  command.
0      1    *   *     *     /Applications/MAMP/htdocs/ztp/zentaopms/bin/dailyreminder.sh   # dailyreminder.
1      1    *   *     *     /Applications/MAMP/htdocs/ztp/zentaopms/bin/backup.sh          # backup database and file.
1      23   *   *     *     /Applications/MAMP/htdocs/ztp/zentaopms/bin/computeburn.sh     # compute burndown chart.
1-59/2 *    *   *     *     /Applications/MAMP/htdocs/ztp/zentaopms/bin/syncsvn.sh         # sync subversion.
1-59/2 *    *   *     *     /Applications/MAMP/htdocs/ztp/zentaopms/bin/syncgit.sh         # sync git
